# extra/webtagger.service
[Unit]
Description=Webtagging service

[Service]
Type=simple
ExecStart=/usr/bin/webtagger -d -c /etc/webtagger/

[Install]
WantedBy=multi-user.target